# install/com.whisper.stream-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.whisper.stream-server</string>

    <!-- Program and arguments -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/whisper-stream-server</string>
        <string>--model</string>
        <string>/usr/local/share/whisper/ggml-base.en.bin</string>
        <string>--port</string>
        <string>9090</string>
        <string>--contexts</string>
        <string>2</string>
        <string>--threads</string>
        <string>4</string>
    </array>

    <!-- Start on boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on failure -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Restart throttling (wait 10 seconds between restarts) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/usr/local/share/whisper</string>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/whisper-stream-server.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/whisper-stream-server.error.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Resource limits (optional) -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>

    <!-- Process priority (optional, lower = higher priority) -->
    <key>Nice</key>
    <integer>-5</integer>
</dict>
</plist>
